# paired/live-keeper/systemd/bili-live-keeper-push.service
[Unit]
Description=Bilibili Live Keeper FFmpeg Keepalive Push
After=network-online.target bili-live-keeper.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bili-live-keeper
EnvironmentFile=/opt/bili-live-keeper/.env
Environment=HOME=/var/lib/bili-live-keeper
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/bili-live-keeper/.venv/bin/python -m bili_live_keeper.cli push-keepalive
Restart=always
RestartSec=10
User=bili-live-keeper
Group=bili-live-keeper
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/bili-live-keeper /var/log/bili-live-keeper /opt/bili-live-keeper/runtime /opt/bili-live-keeper/logs
CapabilityBoundingSet=
LockPersonality=true
UMask=0077

[Install]
WantedBy=multi-user.target
